# systemd/crypto-retention.timer
[Unit]
Description=Daily timer for Crypto Toolkit retention manager

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target